# ha_forge/rootfs/etc/services.d/ha_forge/run
#!/usr/bin/with-contenv sh
echo "Starting HA Forge add-on..."
echo "Node: $(node --version)"
exec node /app/packages/addon/dist/index.js
